# com.dovos.rag.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dovos.rag</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/markrichman/projects/dovos/start_production.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/markrichman/projects/dovos</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/markrichman/projects/dovos/logs/dovos-rag.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/markrichman/projects/dovos/logs/dovos-rag.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
